# scripts/helm.service
[Unit]
Description=Helm API + Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/helm
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/helm/start.sh
ExecStop=%h/helm/stop.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
