# scripts/switchboard.service
[Unit]
Description=Switchboard API + Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/switchboard
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/switchboard/start.sh
ExecStop=%h/switchboard/stop.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
