# contrib/telegram/helm-telegram.service
[Unit]
Description=Helm Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/helm/contrib/telegram/bot
ExecStart=/usr/bin/python3 %h/helm/contrib/telegram/bot/telegram_bot.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
